# raspberry-pi-3/Boot_launch/acces_point.service
[Unit]
description=lance le hotspot

[Service]
ExecStart=/usr/local/bin/acces_point.sh

[Install]
WantedBy=multi-user.target

