# raspberry-pi-3/Boot_launch/web_server.service
[Unit]
Description=Lancer serveur web

[Service]
ExecStart=/usr/local/bin/web_server.sh

[Install]
WantedBy=multi-user.target
